# install/nc/backup/restic-backup.timer
[Timer]
OnCalendar=daily
AccuracySec=5h
Unit=restic-backup.service

[Install]
WantedBy=timers.target
